# scripts/aws-refresh-token.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.user.awsrefreshtoken</string>

    <key>ProgramArguments</key>
    <array>
      <string>/Users/jimweller/.config/dotfiles/scripts/aws-refresh-token.sh</string>
    </array>

    <key>RunAtLoad</key>
    <true />

    <!-- Runs at 00:00, 09:00, and 18:00 -->
    <key>StartCalendarInterval</key>
    <array>
      <dict>
        <key>Hour</key>
        <integer>0</integer>
        <key>Minute</key>
        <integer>0</integer>
      </dict>
      <dict>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
      </dict>
      <dict>
        <key>Hour</key>
        <integer>18</integer>
        <key>Minute</key>
        <integer>0</integer>
      </dict>
    </array>

    <key>StandardOutPath</key>
    <string>/Users/jimweller/assets/aws/launchd.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/jimweller/assets/aws/launchd.err</string>
  </dict>
</plist>