# scripts/sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Nightly encrypted backup to Google Drive sparse image
  
  Install:
    cp scripts/sync.plist ~/Library/LaunchAgents/com.user.sync.plist
    launchctl load ~/Library/LaunchAgents/com.user.sync.plist
  
  Uninstall:
    launchctl unload ~/Library/LaunchAgents/com.user.sync.plist
    rm ~/Library/LaunchAgents/com.user.sync.plist
  
  Start immediately:
    launchctl start com.user.sync
  
  Check status:
    launchctl list | grep sync
-->
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.user.sync</string>

    <key>ProgramArguments</key>
    <array>
      <string>/Users/jimweller/.config/dotfiles/scripts/sync.sh</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
      <key>PATH</key>
      <string>/Users/jimweller/.asdf/shims:/Users/jimweller/.go/bin:/Users/jimweller/.config/dotfiles/scripts:/Users/jimweller/.local/bin:/Users/jimweller/bin:/nix/var/nix/profiles/default/bin:/opt/homebrew/opt/curl/bin:/opt/homebrew/bin:/opt/homebrew/sbin:/Users/jimweller/go/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
      <key>HOME</key>
      <string>/Users/jimweller</string>
      <key>GOPATH</key>
      <string>/Users/jimweller/.go</string>
      <key>ASDF_DATA_DIR</key>
      <string>/Users/jimweller/.asdf</string>
    </dict>

    <key>RunAtLoad</key>
    <true />

    <!-- Runs daily at 2:00 AM -->
    <key>StartCalendarInterval</key>
    <dict>
      <key>Hour</key>
      <integer>2</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/jimweller/.logs/backup-log.txt</string>

    <key>StandardErrorPath</key>
    <string>/Users/jimweller/.logs/backup-err.txt</string>
  </dict>
</plist>